# contrib/systemd/riku-nginx.path
[Unit]
Description=Watch for Riku nginx configuration changes
Documentation=https://dreygur.github.io/riku/
PartOf=riku.service

[Path]
PathModified=/home/deploy/.riku/nginx
Unit=riku-nginx-reload.service

[Install]
WantedBy=multi-user.target
